# service/inetbox2mqtt/run
#!/bin/sh
exec 2>&1
exec /usr/bin/python3 /data/inetbox2mqtt/src/main.py
